# src/systemd/mpd-unload-cdrom.service
[Service]
Type=oneshot
ExecStart=/usr/bin/mpd_cdrom unload
